# supervisor-confd.conf
[program:confd]
command=/data/gopath/bin/confd
autostart=true
autorestart=true
startsecs=5
